# ansible/roles/pihole-maintenance/templates/pihole-update-notify.service.j2
[Unit]
Description=Send Telegram notification on pihole-update failure

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -fsS -o /dev/null -X POST \
  https://api.telegram.org/bot{{ vault_telegram_bot_token }}/sendMessage \
  --data-urlencode "chat_id={{ vault_telegram_chat_id }}" \
  --data-urlencode "text=⚠ pihole-update failed on {{ inventory_hostname }} ({{ ansible_default_ipv4.address | default('?') }}) — check journalctl -u {{ pihole_maint_service_name }}.service"
