# ansible/roles/sealed-secrets-backup/templates/sealed-secrets-backup.service.j2
[Unit]
Description=Encrypted backup of Sealed Secrets master key (age + NFS)
Documentation=file://{{ sealed_backup_script_path }}
After=network-online.target snap.microk8s.daemon-kubelite.service
Wants=network-online.target
OnFailure={{ sealed_backup_service_name }}-notify.service

[Service]
Type=oneshot
User=root
ExecStart={{ sealed_backup_script_path }}

# Light hardening — NoNewPrivileges removed because snap microk8s.kubectl
# uses setuid to transition into its AppArmor profile, which is blocked
# by NoNewPrivileges=true (see "cannot change profile for the next exec call").
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
