# ansible/roles/desktop/templates/cpu-governor.service.j2
# Managed by Ansible (desktop role). Force CPU governor at boot.
# i3-5005U (intel_pstate) bride en powersave par défaut -> {{ desktop_cpu_governor }}
# pour la vidéo/jeu. Pas de paquet (cpupower) requis : echo direct dans sysfs.
[Unit]
Description=Set CPU scaling governor to {{ desktop_cpu_governor }}
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for g in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do echo {{ desktop_cpu_governor }} > "$g"; done'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
