# ansible/roles/desktop/templates/cpu-governor.service.j2
# Managed by Ansible (desktop role). Tuning CPU au boot.
# Intel : force le gouverneur (bride en powersave sinon).
# AMD Zen3 (amd_pstate=active) : laisser le gouverneur dynamique, poser l'EPP.
[Unit]
Description=CPU tuning (governor + EPP)
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c '{% if desktop_force_cpu_governor | bool %}for g in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do echo {{ desktop_cpu_governor }} > "$g"; done; {% endif %}{% if desktop_cpu_epp | length > 0 %}for e in /sys/devices/system/cpu/cpu*/cpufreq/energy_performance_preference; do echo {{ desktop_cpu_epp }} > "$e" 2>/dev/null || true; done; {% endif %}true'

[Install]
WantedBy=multi-user.target
